# killdisk.service
[Unit]
Description=killdisks

[Service]
Type=simple
ExecStart=/usr/bin/killdisk

[Install]
WantedBy=multi-user.target
